# roles/usbmidihost/templates/usbmidihost.service.j2
[Unit]
Description=Initial USB MIDI connect

[Service]
ExecStartPre=/bin/sleep 30
ExecStart={{ usbmidihost_ExecStart }}

[Install]
WantedBy=multi-user.target
